# pulse_audio_stereo_pairs/root/etc/s6-overlay/s6-rc.d/svc-pulse_audio_stereo_pairs/finish
#!/bin/bash
LOG="/config/local_pulse_audio_stereo_pairs/remap_modules.log"
echo "$(date) — addon stopping, unloading remap sinks..." >> "$LOG"
for id in $(pactl list short modules | awk '/module-remap-sink/ {print $1}'); do
    pactl unload-module "$id"
done
echo "$(date) — remap sinks unloaded" >> "$LOG"